# etc/systemd/system/titdb.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/titdb -d /dev/input/by-path/platform-39b10c000.spi-cs-0-event-mouse -m f
